# ansible/roles/ml-server/templates/headroom.service.j2
[Unit]
Description=Headroom AI Proxy
After=network.target
Documentation=https://headroom.ai

[Service]
Type=simple
ExecStart={{ ml_install_prefix }}/miniforge3/envs/venv-system/bin/headroom proxy
User=headroom
Group=headroom
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=headroom

Environment="PATH={{ ml_install_prefix }}/miniforge3/envs/venv-system/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"

[Install]
WantedBy=multi-user.target
